# filesystem/.config/autostart/Cryptomator.desktop
[Desktop Entry]
Type=Application
Exec=cryptomator
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Cryptomator
Comment=Created with org.cryptomator.linux.autostart.FreedesktopAutoStartService
