# filesystem/.config/autostart/discord.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/discord --tray --start-minimized
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=Discord
Name=Discord
Comment[en_US]=Discord text and video chat
Comment=Discord text and video chat
